# deploy/systemd/smart-db-backup.service
[Unit]
Description=Smart DB backup job
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
WorkingDirectory=/opt/smart-db/deploy
ExecStart=/opt/smart-db/deploy/scripts/backup.sh
